# conf.d/fluxbox.conf
[program:fluxbox]
command=fluxbox %(ENV_WM_OPTIONS)s
autostart=%(ENV_AUTO_START_WM)s
autorestart=true
# stdout_logfile=/var/log/fluxbox.log
# stderr_logfile=/var/log/fluxbox.err.log
